# cosmic-core/systemvm/patches/centos7/etc/udev/rules.d/75-persistent-net-generator.rules
# Bring up any new plugged interface
ACTION=="add", SUBSYSTEM=="net",  RUN+="/usr/sbin/ip link set dev '%k' up"
